# contrib/lazyfwd.service
[Unit]
Description=lazyfwd - On-demand Kubernetes port forwarding
Documentation=https://github.com/atas/lazyfwd
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/lazyfwd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
